# server/supervisor.conf
[program:MYDOMAIN]
command=MYGUNICORN
directory=MYDIR
autostart=true
autorestart=true
user=MYUSER
stdout_logfile=MYLOG_ACCESS
stderr_logfile=MYLOG_ERROR